# templates/FMEServerCleanupStart.erb
#!/bin/sh
# chkconfig: 3 97 99
### BEGIN INIT INFO
# Provides:           FMEServerCleanupStart
# Required-Start:     $local_fs
# Required-Stop:      $local_fs
# Default-Start: 3
# Default-Stop: 6
# Short-Description:  FME Server Cleanup Start
# Description:        FME Server Cleanup Start
### END INIT INFO

RETVAL=0

# start, debug, stop, and status functions
start() {
        <%= @install_directory %>/Server/startCleanup.sh &
        RETVAL=$?
}

stop() {
        <%= @install_directory %>/Server/stopCleanup.sh
        RETVAL=$?
}

status() {
        ps -ef | grep -v grep | grep '<%= @install_directory %>/Server/startCleanup.sh' > /dev/null 2>&1
        RETVAL=$?
        if [ $RETVAL -ne 0 ]; then
                echo "FMEServerCleanupStart stopped"
        else
                echo "FMEServerCleanupStart running"
        fi
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  status)
        status
        ;;
  restart)
        stop
        start
        ;;
  *)
        echo "Usage: $0 {start|stop|status|restart}"
        exit 1
esac

exit $RETVAL
