# templates/FMEServerDatabaseStart.erb
#!/bin/sh
# chkconfig: 3 98 99
### BEGIN INIT INFO
# Provides:           FMEServerDatabaseServerStart
# Required-Start:     $local_fs
# Required-Stop:      $local_fs
# Default-Start: 3
# Default-Stop: 6
# Short-Description:  FME Database Server Start
# Description:        FME Database Server Start
### END INIT INFO

start() {
    su <%= @user %> <%= @install_directory %>/Server/startDatabaseServer.sh
    RETVAL=$?
}

stop() {
    su <%= @user %> <%= @install_directory %>/Server/stopDatabaseServer.sh
    RETVAL=$?
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  *)
    echo "Usage: $0 {start|stop}"
    exit 1
esac

exit $RETVAL
