# templates/FMEServerAppServerStart.erb
#!/bin/sh
# chkconfig: 3 97 99
### BEGIN INIT INFO
# Provides:           FMEServerAppserverStart
# Required-Start:     $local_fs
# Required-Stop:      $local_fs
# Default-Start: 3
# Default-Stop: 6
# Short-Description:  FME Server Start
# Description:        FME Server Start
### END INIT INFO

JAVA_HOME=<%= @install_directory %>/Utilities/jre
export JAVA_HOME

RETVAL=0

# start, debug, stop, and status functions
start() {
    <%= @install_directory %>/Server/startApplicationServer.sh
    RETVAL=$?
}

stop() {
    <%= @install_directory %>/Server/stopApplicationServer.sh
    RETVAL=$?
}

status() {
    ps -p `cat /var/run/catalina.pid` > /dev/null 2>&1
    RETVAL=$?
    if [ $RETVAL -ne 0 ]; then
        echo "FMEServerAppserverStart stopped"
    else
        echo "FMEServerAppserverStart running"
    fi
}

restart() {
    stop
    start
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  status)
        status
        ;;
  restart)
        restart
        ;;
  *)
    echo "Usage: $0 {start|stop|status|restart}"
    exit 1
esac

exit $RETVAL
