# templates/FMEServerSMTPRelayStart.erb
#!/bin/sh
# chkconfig: 3 98 99
### BEGIN INIT INFO
# Provides:          FMEServerSMTPRelay
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     3
# Default-Stop:      6
# Short-Description: Startup of FME Server SMTP Relay
# Description: Startup of FME Server SMTP Relay
### END INIT INFO

RETVAL=0

# start, debug, stop, and status functions
start() {
        <%= @install_directory %>/Server/startFMEServerSMTPRelay.sh &
        RETVAL=$?
}

stop() {
        <%= @install_directory %>/Server/stopFMEServerSMTPRelay.sh
        RETVAL=$?
}

status() {
        JAVA_HOME=<%= @install_directory %>/Utilities/jre <%= @install_directory %>/Utilities/smtprelay/james/bin/phoenix.sh check > /dev/null 2>&1
        RETVAL=$?
        if [ $RETVAL -ne 0 ]; then
                echo "FMEServerSMTPRelay stopped"
        else
                echo "FMEServerSMTPRelay running"
        fi
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;  
  status)
        status
        ;;
  restart)
        stop
        start
        ;;
  *)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
esac

exit $RETVAL
